# tg-secretary.service
[Unit]
   Description=Telegram Secretary Bot
   After=network.target

   [Service]
   Type=simple
   User=your_username
   Group=your_username
   WorkingDirectory=/path/to/TG-secretary
   Environment=PYTHONPATH=/path/to/TG-secretary
   Environment=PATH=/path/to/TG-secretary/.venv/bin:/Users/kristina/Desktop/Denis/TG-secretary/.venv/bin:/Users/kristina/.pyenv/shims:/Users/kristina/.pyenv/bin:/opt/homebrew/bin:/opt/homebrew/sbin:/Library/Frameworks/Python.framework/Versions/3.13/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin:/Users/kristina/.local/bin
   ExecStart=/path/to/TG-secretary/.venv/bin/python main.py
   Restart=always
   RestartSec=10

   [Install]
   WantedBy=multi-user.target
